# templates/debian-buster-mate/home/jitas/.config/autostart/mate-power-manager.desktop
[Desktop Entry]
Type=Application
Name=Power Manager
Exec=mate-power-manager
Terminal=false
OnlyShowIn=MATE;
X-MATE-Autostart-enabled=false
