# templates/debian-bullseye-mate/etc/skel/.config/systemd/user/video.service
[Unit]
Description=video service

[Service]
ExecStart=/usr/local/bin/video-start
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
