# templates/debian-bullseye-mate/etc/skel/.config/autostart/mate-screensaver.desktop
[Desktop Entry]
Type=Application
Name=Screensaver
Exec=mate-screensaver
OnlyShowIn=MATE;
NoDisplay=true
X-MATE-Autostart-enabled=false
