# templates/debian-bullseye-mate/etc/skel/.config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network
Exec=nm-applet
Terminal=false
NoDisplay=true
NotShowIn=KDE;GNOME;
X-MATE-Autostart-enabled=false
