# templates/debian-buster-mate/etc/systemd/system/certbot.service.d/override.conf
[Service]
ExecStartPost=
ExecStartPost=/bin/systemctl restart websockify-secure.service
ExecStartPost=/bin/systemctl reload nginx.service
